# systemd/hermes-ha-monitor.service
[Unit]
Description=Publish Hermes host and Codex metrics via MQTT
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/hermes-ha-monitor.env
ExecStart=/opt/hermes-ha-monitor/hermes_ha_monitor.py
User=root
Group=root
Nice=5
IOSchedulingClass=best-effort
ProtectSystem=full
ReadWritePaths=/var/lib/hermes-ha-monitor /root/.cache /root/.codex/log
PrivateTmp=true
